# debian/wazo-ui.service
[Unit]
Description=wazo-ui
After=network.target
Before=monit.service

[Service]
Type=forking
ExecStartPre=/usr/bin/install -d -o wazo-ui -g wazo-ui /run/wazo-ui
ExecStart=/usr/bin/wazo-ui
PIDFile=/run/wazo-ui/wazo-ui.pid

[Install]
WantedBy=multi-user.target
